# systemd/cynara-admin.socket.in
[Socket]
ListenStream=/run/cynara/cynara-admin.socket
SocketMode=0700
SmackLabelIPIn=@
SmackLabelIPOut=@

Service=cynara.service

[Unit]
@SOCKET_UNIT_SETTINGS@

[Install]
WantedBy=sockets.target
