# systemd/cynara.socket.in
[Socket]
ListenStream=/run/cynara/cynara.socket
SocketMode=0777
SmackLabelIPIn=*
SmackLabelIPOut=@

Service=cynara.service

[Unit]
@SOCKET_UNIT_SETTINGS@

[Install]
WantedBy=sockets.target
