# systemd/cynara-monitor-get.socket.in
[Socket]
ListenStream=/run/cynara/cynara-monitor-get.socket
SocketGroup=security_fw
SocketMode=0060
SmackLabelIPIn=@
SmackLabelIPOut=@

Service=cynara.service

[Unit]
@SOCKET_UNIT_SETTINGS@

[Install]
WantedBy=sockets.target
